# deploy/ecommerce.service
[Unit]
Description=E-commerce Platform
After=network.target postgresql.service docker.service

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/E-commerce-Platform
Environment=NODE_ENV=production
EnvironmentFile=/home/azureuser/E-commerce-Platform/.env
ExecStart=/usr/bin/node /home/azureuser/E-commerce-Platform/dist/index.cjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
